# com.chiha.youtube.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chiha.youtube.daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/chiha/projects-active/youtube/scripts/run_daily_youtube_tasks.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/chiha/projects-active/youtube</string>
    <key>StandardOutPath</key>
    <string>/Users/chiha/projects-active/youtube/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chiha/projects-active/youtube/launchd_stderr.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>Persistent</key>
    <true/>
</dict>
</plist>
